# conf/cyndikator.service
[Unit]
Description=Cyndikator daemon
[Service]
Type=simple
ExecStart=${HOME}/.cargo/bin/cynd run
[Install]
WantedBy=default.target

